# consumers/consumer.conf
[program:KFCProviderSearchAtrapaloUsers]
command=python /consumers/consumer.py
autostart=true
autorestart=true
numprocs=1
process_name = %(program_name)s_%(process_num)02d
startsecs=0